# coflnet-start.service
[Unit]
Description=Runs on startup, Turn swap off 

[Service]
ExecStart=/sbin/swapoff -a

[Install]
WantedBy=multi-user.target
